# debian/python3-cyberfusion-rabbitmq-consumer-log-server.cron.d
5 4 * * * root /usr/bin/tuxis-cron /etc/cronscripts/rabbitmq-consumer-log-server-purge-logs.cron
